# infra/systemd/guardrail-agent.service
[Unit]
Description=Guardrail Alpha Agent

[Service]
WorkingDirectory=/opt/guardrail-alpha
ExecStart=/usr/local/bin/guardrail-agent --config configs/production.toml
Restart=always

[Install]
WantedBy=multi-user.target

